# deploy/lumenalpha-refresh.timer
[Unit]
Description=Run LumenAlpha daily stock analysis refresh at 19:30 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 19:30:00
Persistent=true
Unit=lumenalpha-refresh.service

[Install]
WantedBy=timers.target
