# deploy/lumenalpha-refresh.timer
[Unit]
Description=Run LumenAlpha daily stock analysis refresh at 21:00 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 21:00:00
Persistent=true
Unit=lumenalpha-refresh.service

[Install]
WantedBy=timers.target
